# templates/etc/cron.daily/gitlab.j2
#!/bin/sh
/opt/gitlab/bin/gitlab-backup create
